# layers/meta-balena-jetson/recipes-bsp/tegra-binaries/files/nexcom.service
[Unit]
Description=The OS status Led Service

[Service]
Type=simple
ExecStart=/etc/systemd/ledstate.sh boot

[Install]
WantedBy=sysinit.target
