# layers/meta-balena-jetson/recipes-bsp/tegra-binaries/files/usb_power.service
[Unit]
Description=Nexcom control usb power
DefaultDependencies=no
After=systemd-reboot.service

[Service]
Type=simple
ExecStart=/etc/systemd/ledstate.sh usb-resume

[Install]
WantedBy=reboot.target
